# rpm/penstock.service
[Unit]
Description=Penstock
After=network.target

[Service]
Type=simple
User=psale
ExecStart=/opt/psale/penstock/bin/penstock /etc/penstock/penstock.yml
Restart=on-abort

[Install]
WantedBy=multi-user.target
